# kikiriki.service
[Unit]
Description=Kikiriki Speech-to-Text Push-to-Talk
Documentation=https://github.com/rlaffers/kikiriki
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=KIKIRIKI_INSTALL_DIR/kikiriki-bind-keys
ExecStop=KIKIRIKI_INSTALL_DIR/kikiriki-unbind-keys
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Kill all child processes when stopping
KillMode=control-group

[Install]
WantedBy=graphical-session.target
